# infra/linkedin-scraper.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd schedule for the LinkedIn scraper.

  IMPORTANT — launchd does NOT wake a sleeping Mac. StartCalendarInterval only
  runs the job when the machine is already awake (or on its next wake). To wake
  the Mac for the overnight run, install the daily wake ONCE:
    sudo ./scripts/setup_wake_schedule.sh      # pmset repeat wake @ 21:58
  The wrapper then runs the scrape under `caffeinate` (stays awake through the
  upload) and, because of the --sleep-after arg below, sleeps the Mac when done.

  Install:
    cp infra/linkedin-scraper.plist.example ~/Library/LaunchAgents/com.<you>.linkedin-scraper.plist
    # edit the two /path/to/ds_jobs paths below, then:
    launchctl load ~/Library/LaunchAgents/com.<you>.linkedin-scraper.plist
-->
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.user.linkedin-scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/ds_jobs/scripts/run_daily.sh</string>
        <string>--sleep-after</string>
    </array>

    <!-- Fire daily at 10:00 pm. The Mac is woken at 21:58 by the pmset repeat
         wake from setup_wake_schedule.sh (run that once); this only sets when
         launchd runs the job, it does not wake the machine on its own. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/path/to/ds_jobs/logs/scrape.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/ds_jobs/logs/scrape.log</string>

</dict>
</plist>
